# docker/etc/cron.d/cloudsummariser
10 2,14 * * * root /usr/bin/run_cloud_summariser.sh
